# supervisord.conf
[supervisord]
logfile=%(here)s/log/supervisord.log
pidfile=%(here)s/log/supervisord.pid


[unix_http_server]
file = /tmp/triagesup.sock


[supervisorctl]
serverurl = unix:///tmp/triagesup.sock
prompt = triage


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:api]
command=%(here)s/env/bin/python %(here)s/triage/api.py %(here)s/development.ini
redirect_stderr=true
stdout_logfile=%(here)s/log/api.log


[program:web]
command=%(here)s/env/bin/pserve %(here)s/development.ini --reload --monitor-restart
stdout_logfile=%(here)s/log/web.log
stderr_logfile=%(here)s/log/weberrors.log
